# rootfs/etc/s6/services/cron/run
#!/bin/sh
cd /var/www/html
while true; do
  php artisan schedule:run
  sleep 60
done
