# rootfs/etc/s6/services/queue/run
#!/bin/sh
cd /var/www/html || exit 1
exec php artisan queue:work --queue=emails,default --sleep=3 --tries=3
